# templates/zookeeper.service.erb
# Autogenerated by Puppet from zookeeper/templates/zookeeper.service.erb
# Do not edit

[Unit]
Description=Apache ZooKeeper
<%if @systemd_unit_want -%>
Want=<%=@systemd_unit_want %>
<% end -%>
<%if @systemd_unit_after -%>
After=<%=@systemd_unit_after %>
<% end -%>

[Service]
ExecStart=/bin/sh -c 'set -x; . <%= @cfg_dir %>/environment; CLASSPATH="<%= @zoo_dir %>/zookeeper.jar:<%= @zoo_dir %>/lib/*:$CLASSPATH"; CLASSPATH="$(. <%= @zoo_dir %>/bin/zkEnv.sh ; echo $CLASSPATH)"; mkdir -p <%= @log_dir %>; $JAVA "-Dzookeeper.log.dir=<%= @log_dir %>" "-Dzookeeper.root.logger=$ZOO_LOG4J_PROP" -cp "$CLASSPATH" $JAVA_OPTS "$ZOOMAIN" "$ZOOCFG"'
Restart=always
RemainAfterExit=no
SyslogIdentifier=zookeeper
User=<%= @user %>
Group=<%= @group %>

[Install]
WantedBy=multi-user.target
